# debian/eggnoggpi/DEBIAN/postinst
#!/bin/bash

echo "Install the service"

systemctl daemon-reload
systemctl enable mpv eggnoggpi pulseaudio
systemctl start mpv eggnoggpi pulseaudio
systemctl set-default graphical.target

echo "@reboot    root    mkdir -p /run/user/1000/pulse; chown -R pi:pi /run/user/1000" >> /etc/crontab
