# debian/eggnoggpi/DEBIAN/postinst
#!/bin/bash

echo "Install the service"

systemctl daemon-reload
systemctl enable mpv nodm eggnoggpi
systemctl start mpv nodm eggnoggpi
echo "@reboot    root    sleep 15; /sbin/service eggnoggpi start" >> /etc/crontab
