# debian/eggnoggpi/DEBIAN/preinst
#!/bin/bash

echo "Setup swap to no ..."
dphys-swapfile swapoff
dphys-swapfile uninstall
systemctl disable dphys-swapfile

