# debian/eggnoggpi/DEBIAN/prerm
#!/bin/bash

echo "disabling services"
systemctl disable eggnoggpi mpv antimicrox
systemctl stop eggnoggpi mpv antimicrox
